# templates/default/chef-zero.conf.erb
description "Chef-Zero"

start on runlevel [2345]
stop on shutdown
respawn

script
  exec >> /var/log/$UPSTART_JOB.log 2>&1
  set -x
  # start chef-zero service
  exec sudo -i -u root <%= @chef_embed_path %>/chef-zero -p <%= @chef_server_port %> -l info
end script
